# com.apple-assistant.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.apple-assistant</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/I306141/Documents/Automations/apple_assistant.sh</string>
    </array>

    <key>WatchPaths</key>
    <array>
        <!-- Apple Notes SQLite database -->
        <string>/Users/I306141/Library/Group Containers/group.com.apple.notes/NoteStore.sqlite</string>
        <!-- Apple Reminders calendar store -->
        <string>/Users/I306141/Library/Calendars</string>
    </array>

    <key>ThrottleInterval</key>
    <integer>15</integer>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/I306141/Documents/Automations/launchd_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/I306141/Documents/Automations/launchd_stderr.log</string>
</dict>
</plist>
